# allwinner/sayeye/sayeye.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=20
STOP=20

USE_PROCD=1
PROG=/bin/sayeye
OOM_ADJ=-17

start_service() {
	mkdir -p /dev/socket

	procd_open_instance
	procd_set_param oom_adj $OOM_ADJ
	procd_set_param command $PROG -D
	procd_close_instance
}

shutdown() {
	echo "sayeye exit"
}
